# eidos_mcp/systemd/eidos-mcp.service
[Unit]
Description=Eidosian MCP Server (StreamableHTTP)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/lloyd/eidosian_forge/eidos_mcp
Environment=PYTHONPATH=/home/lloyd/eidosian_forge/eidos_mcp/src
Environment=EIDOS_FORGE_DIR=/home/lloyd/eidosian_forge
Environment=EIDOS_MCP_TRANSPORT=streamable-http
Environment=EIDOS_MCP_MOUNT_PATH=/mcp
Environment=EIDOS_MCP_STATELESS_HTTP=1
Environment=EIDOS_MCP_ENABLE_COMPAT_HEADERS=1
Environment=EIDOS_MCP_ENABLE_SESSION_RECOVERY=1
Environment=EIDOS_MCP_ENABLE_ERROR_RESPONSE_COMPAT=1
Environment=FASTMCP_HOST=127.0.0.1
Environment=FASTMCP_PORT=8928
Environment=EIDOS_OAUTH2_PROVIDER=google
Environment=EIDOS_OAUTH2_AUDIENCE=
Environment=EIDOS_OAUTH2_STATIC_BEARER=
Environment=EIDOS_MCP_SERVICE_STATE=/home/lloyd/.eidosian/run/mcp_service_state.json
Environment=EIDOS_MCP_BACKUP_DIR=/home/lloyd/.eidosian/backups/eidos_mcp
ExecStart=/home/lloyd/eidosian_forge/eidosian_venv/bin/python3 -c "import eidos_mcp.eidos_mcp_server as s; s.main()"
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
